# serverroom/systemd/mass-db-writer.service
[Unit]
Description=MASS Database Writer
After=mass-startup.service docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/sleep 45
ExecStart=/bin/bash -c "docker exec campus-postgres psql -U massadmin -c 'CREATE DATABASE massdb;' postgres 2>/dev/null || true && docker run --rm -d --name mass-db-writer --network campus-mass-network --no-healthcheck -v /home/menna/agents-test-main/agents-test-main/db_writer.py:/app/db_writer.py docker-ndr-agent python3 /app/db_writer.py"
ExecStop=/bin/bash -c "docker stop mass-db-writer 2>/dev/null || true"

[Install]
WantedBy=multi-user.target
